# systemd/calendar.service
[Unit]
Description=Pi Family Calendar Web App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__/app
EnvironmentFile=-/boot/firmware/calendar.env
Environment=NODE_ENV=production
Environment=CALENDAR_DIR=/home/__USER__/.local/share/calendar
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
