# systemd/kiosk-reboot.timer
[Unit]
Description=Run kiosk reboot escalation every 5 minutes
After=graphical.target

[Timer]
OnBootSec=15min
OnUnitActiveSec=5min
Unit=kiosk-reboot.service

[Install]
WantedBy=timers.target
